# common/files/ansible.cron
#!/bin/bash

# Run ansible

ansible-playbook -c local /etc/ansible/src/site.yml