# elasticluster/share/playbooks/roles/mesos-common/files/etc/profile.d/pyspark.sh
#!/usr/bin/env bash
# Additional settings for PySpark
# by default, calling `pyspark` starts an IPython interpreter
export PYSPARK_DRIVER_PYTHON=ipython
