# elasticluster/share/playbooks/roles/pdsh/files/etc/profile.d/pdsh.sh
# THIS FILE IS CONTROLLED BY ELASTICLUSTER
# local modifications will be overwritten
# the next time `elasticluster setup` is run!
#

# PDSH defaults to using `rsh` for executing remote commands, but SSH
# is actually a much better default choice...
export PDSH_RCMD_TYPE='ssh'
