# elasticluster/share/playbooks/roles/spark-common/files/etc/profile.d/alluxio.sh
#!/usr/bin/env bash
export PATH=$PATH:$ALLUXIO_HOME/bin