# elasticluster/share/playbooks/roles/spark-common/files/etc/profile.d/spark.sh
#!/usr/bin/env bash
export PATH=$PATH:$SPARK_HOME/bin