# elasticluster/share/playbooks/roles/jupyterhub/files/etc/supervisor/conf.d/jupyterhub.conf
# THIS FILE IS CONTROLLED BY ELASTICLUSTER
# local modifications will be overwritten
# the next time `elasticluster setup` is run!
#

# Control (auto)start of JupyterHub through `supervisor` (http://supervisord.org/)
#
# Initial version taken from:
# https://github.com/jupyterhub/jupyterhub-tutorial/blob/master/supervisor/jupyterhub.conf

[program:jupyterhub]
command=/opt/anaconda3/bin/jupyterhub -f /etc/jupyterhub/jupyterhub_config.py
directory=/var/lib/jupyterhub
autostart=true
autorestart=true
startretries=3
exitcodes=0,2
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/log/jupyterhub.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
user=root
